# host-conf/upstart/logstash.conf
description "Logstash container"
author "klynch@gmail.com"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start logstash
end script
